# dotfiles/.config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager Applet

[Service]
Type=simple
ExecStart=nm-applet

[Install]
WantedBy=default.target
