# dotfiles/.config/systemd/user/protonmail-bridge.service
[Unit]
Description=IMAP/POP3 bridge for ProtonMail

[Service]
Type=simple
ExecStart=/usr/bin/protonmail-bridge

[Install]
WantedBy=default.target
